# crontab
# mm hh jj MMM JJJ tâche > log
00 02   *  *    *   /usr/bin/php5 /data/services/web/www.ethic-etapes.fr/admin/specifique/sejourDevalidation.action.php
